# test/server/api_test_redmine.service
[Unit]
Description=api_test_redmine

[Service]
ExecStart=/home/test/api/api_test_redmine.sh

[Install]
WantedBy=multi-user.target 