# test/server/api_test_alfresco.service
[Unit]
Description=api_test_alfresco

[Service]
ExecStart=/home/test/api/api_test_alfresco.sh

[Install]
WantedBy=multi-user.target 